# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info

[program:backend]
; FastAPI via Uvicorn (passt zu backend/main.py)
command=/usr/local/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/backend_stdout.log
stderr_logfile=/var/log/backend_stderr.log

[program:wstation]
; Witterungsstation-Anwendung
command=python -u witterungsstation_py/main.py
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/wstation_stdout.log
stderr_logfile=/var/log/wstation_stderr.log

[program:frontend_static]
; Einfacher statischer Server für gebuildetes Frontend
; Serviert /app/frontend_dist auf Port 80
command=python -u -m http.server 80 --directory /app/frontend_dist
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/frontend_stdout.log
stderr_logfile=/var/log/frontend_stderr.log
